# _docker/app/crontab
* * * * * cd /var/www/chat && /usr/local/bin/php artisan set-offline-status >> /var/log/cron.log 2>&1
